# deploy/epson-rp-api.service
[Unit]
Description=Epson TM-T88VI receipt printer HTTP API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=plugdev
WorkingDirectory=/home/pi/epson-rp-api
EnvironmentFile=-/etc/default/epson-rp-api
ExecStart=/home/pi/epson-rp-api/venv/bin/python -m printapi
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
